# .config/autostart/polybar.desktop
[Desktop Entry]
Type=Application
Exec=/home/rile/.config/polybar/launch.sh
Hidden=false
X-GNOME-Autostart-enabled=true
Name=Polybar
